# crontab
0 2 * * * root /usr/local/bin/zfs_snapshot_script.sh
