# archive/python-legacy/crontab
* * * * * root /app/venv/bin/python3 /usr/local/bin/test.py
